# cti-data/indicator/indicator--c125fc72-858b-435a-9cdc-bc6a75712262.json
{"name": "Default PowerSploit and Empire Schtasks Persistence", "description": "Detects the creation of a schtask via PowerSploit or Empire Default Configuration.", "created": "2018-03-06T00:00:00.000Z", "valid_from": "2018-03-06T00:00:00Z", "indicator_types": ["malicious-activity"], "pattern_type": "sigma", "pattern": "title: Default PowerSploit and Empire Schtasks Persistence\nid: 56c217c3-2de2-479b-990f-5c109ba8458f\nstatus: experimental\ndescription: Detects the creation of a schtask via PowerSploit or Empire Default Configuration.\nreferences:\n    - https://github.com/0xdeadbeefJERKY/PowerSploit/blob/8690399ef70d2cad10213575ac67e8fa90ddf7c3/Persistence/Persistence.psm1\n    - https://github.com/EmpireProject/Empire/blob/master/lib/modules/powershell/persistence/userland/schtasks.py\n    - https://github.com/EmpireProject/Empire/blob/master/lib/modules/powershell/persistence/elevated/schtasks.py\nauthor: Markus Neis, @Karneades\ndate: 2018/03/06\nlogsource:\n    product: windows\n    category: process_creation\ndetection:\n    selection:\n        ParentImage:\n            - '*\\powershell.exe'\n        CommandLine:\n            - '*schtasks*/Create*/SC *ONLOGON*/TN *Updater*/TR *powershell*'\n            - '*schtasks*/Create*/SC *DAILY*/TN *Updater*/TR *powershell*'\n            - '*schtasks*/Create*/SC *ONIDLE*/TN *Updater*/TR *powershell*'\n            - '*schtasks*/Create*/SC *Updater*/TN *Updater*/TR *powershell*'\n    condition: selection\ntags:\n    - attack.execution\n    - attack.persistence\n    - attack.privilege_escalation\n    - attack.t1053\n    - attack.t1086\n    - attack.s0111\n    - attack.g0022\n    - attack.g0060\n    - car.2013-08-001\nfalsepositives:\n    - False positives are possible, depends on organisation and processes\nlevel: high\n", "object_marking_refs": ["marking-definition--0a872d70-1d05-45dd-a25f-031af547a102"], "created_by_ref": "identity--082ca4f0-d825-45b6-aa42-ab9afd1e5503", "type": "indicator", "spec_version": "2.1", "id": "indicator--c125fc72-858b-435a-9cdc-bc6a75712262", "modified": "2020-04-02T10:14:44.016Z"}