# cti-data/indicator/indicator--b1cff8ab-02e5-4645-8770-8b471d7efd74.json
{"name": "Rare Schtasks Creations", "description": "Detects rare scheduled tasks creations that only appear a few times per time frame and could reveal password dumpers, backdoor installs or other types of malicious code", "created": "2017-03-23T00:00:00.000Z", "valid_from": "2017-03-23T00:00:00Z", "indicator_types": ["malicious-activity"], "pattern_type": "sigma", "pattern": "title: Rare Schtasks Creations\nid: b0d77106-7bb0-41fe-bd94-d1752164d066\ndescription: Detects rare scheduled tasks creations that only appear a few times per time frame and could reveal password dumpers, backdoor installs or other types of malicious code\nstatus: experimental\nauthor: Florian Roth\ndate: 2017/03/23\ntags:\n    - attack.execution\n    - attack.privilege_escalation\n    - attack.persistence\n    - attack.t1053\n    - car.2013-08-001\nlogsource:\n    product: windows\n    service: security\n    definition: 'The Advanced Audit Policy setting Object Access > Audit Other Object Access Events has to be configured to allow this detection (not in the baseline recommendations by Microsoft). We also recommend extracting the Command field from the embedded XML in the event data.'\ndetection:\n    selection:\n        EventID: 4698\n    timeframe: 7d\n    condition: selection | count() by TaskName < 5\nfalsepositives:\n    - Software installation\n    - Software updates\nlevel: low\n", "object_marking_refs": ["marking-definition--0a872d70-1d05-45dd-a25f-031af547a102"], "created_by_ref": "identity--082ca4f0-d825-45b6-aa42-ab9afd1e5503", "type": "indicator", "spec_version": "2.1", "id": "indicator--b1cff8ab-02e5-4645-8770-8b471d7efd74", "modified": "2020-04-02T10:15:08.998Z"}