# cti-data/indicator--cbe7490e-e7e1-4228-a2f6-29c3e674e4fa.json
{"name": "Persistence and Execution at Scale via GPO Scheduled Task", "description": "Detect lateral movement using GPO scheduled task, ususally used to deploy ransomware at scale", "created": "2019-04-03T00:00:00.000Z", "valid_from": "2019-04-03T00:00:00Z", "indicator_types": ["malicious-activity"], "pattern_type": "sigma", "pattern": "title: Persistence and Execution at Scale via GPO Scheduled Task\nid: a8f29a7b-b137-4446-80a0-b804272f3da2\ndescription: Detect lateral movement using GPO scheduled task, ususally used to deploy ransomware at scale\nauthor: Samir Bousseaden\ndate: 2019/04/03\nreferences:\n    - https://twitter.com/menasec1/status/1106899890377052160\ntags:\n    - attack.persistence\n    - attack.lateral_movement\n    - attack.t1053\nlogsource:\n    product: windows\n    service: security\n    description: 'The advanced audit policy setting \"Object Access > Audit Detailed File Share\" must be configured for Success/Failure'\ndetection:\n    selection:\n        EventID: 5145\n        ShareName: \\\\*\\SYSVOL\n        RelativeTargetName: '*ScheduledTasks.xml'\n        Accesses: '*WriteData*'\n    condition: selection\nfalsepositives:\n    - if the source IP is not localhost then it's super suspicious, better to monitor both local and remote changes to GPO scheduledtasks\nlevel: high\n", "object_marking_refs": ["marking-definition--0a872d70-1d05-45dd-a25f-031af547a102"], "created_by_ref": "identity--082ca4f0-d825-45b6-aa42-ab9afd1e5503", "type": "indicator", "spec_version": "2.1", "id": "indicator--cbe7490e-e7e1-4228-a2f6-29c3e674e4fa", "modified": "2020-04-02T10:15:06.678Z"}