# cti-data/indicator--22d03ecf-c1dc-491d-ac49-d4a8b3cc7077.json
{"name": "Windows 10 Scheduled Task SandboxEscaper 0-day", "description": "Detects Task Scheduler .job import arbitrary DACL write\\par", "created": "2019-05-22T00:00:00.000Z", "valid_from": "2019-05-22T00:00:00Z", "indicator_types": ["malicious-activity"], "pattern_type": "sigma", "pattern": "title: Windows 10 Scheduled Task SandboxEscaper 0-day\nid: 931b6802-d6a6-4267-9ffa-526f57f22aaf\nstatus: experimental\ndescription: Detects Task Scheduler .job import arbitrary DACL write\\par\nreferences:\n    - https://github.com/SandboxEscaper/polarbearrepo/tree/master/bearlpe\nauthor: Olaf Hartong\ndate: 2019/05/22\nlogsource:\n    category: process_creation\n    product: windows\ndetection:\n    selection:\n        Image: schtasks.exe\n        CommandLine: '*/change*/TN*/RU*/RP*'\n    condition: selection\nfalsepositives:\n    - Unknown\ntags:\n    - attack.privilege_escalation\n    - attack.execution\n    - attack.t1053\n    - car.2013-08-001\nlevel: high\n", "object_marking_refs": ["marking-definition--0a872d70-1d05-45dd-a25f-031af547a102"], "created_by_ref": "identity--082ca4f0-d825-45b6-aa42-ab9afd1e5503", "type": "indicator", "spec_version": "2.1", "id": "indicator--22d03ecf-c1dc-491d-ac49-d4a8b3cc7077", "modified": "2020-04-02T10:14:42.702Z"}