# cti-data/indicator--0b994a28-7526-4eac-bbb5-88e3a13980e6.json
{"name": "Scheduled Task Creation", "description": "Detects the creation of scheduled tasks in user session", "created": "2019-01-16T00:00:00.000Z", "valid_from": "2019-01-16T00:00:00Z", "indicator_types": ["malicious-activity"], "pattern_type": "sigma", "pattern": "title: Scheduled Task Creation\nid: 92626ddd-662c-49e3-ac59-f6535f12d189\nstatus: experimental\ndescription: Detects the creation of scheduled tasks in user session\nauthor: Florian Roth\ndate: 2019/01/16\nlogsource:\n    category: process_creation\n    product: windows\ndetection:\n    selection:\n        Image: '*\\schtasks.exe'\n        CommandLine: '* /create *'\n    filter:\n        User: NT AUTHORITY\\SYSTEM\n    condition: selection and not filter\nfields:\n    - CommandLine\n    - ParentCommandLine\ntags:\n    - attack.execution\n    - attack.persistence\n    - attack.privilege_escalation\n    - attack.t1053\n    - attack.s0111\n    - car.2013-08-001\nfalsepositives:\n    - Administrative activity\n    - Software installation\nlevel: low\n", "object_marking_refs": ["marking-definition--0a872d70-1d05-45dd-a25f-031af547a102"], "created_by_ref": "identity--082ca4f0-d825-45b6-aa42-ab9afd1e5503", "type": "indicator", "spec_version": "2.1", "id": "indicator--0b994a28-7526-4eac-bbb5-88e3a13980e6", "modified": "2020-04-02T10:14:46.937Z"}