# cti-data/indicator--4401f7de-c808-42d3-8df1-383e4ce51848.json
{"name": "Rare Scheduled Task Creations", "description": "This rule detects rare scheduled task creations. Typically software gets installed on multiple systems and not only on a few. The aggregation and count function selects tasks with rare names.", "created": "2017-03-17T00:00:00.000Z", "valid_from": "2017-03-17T00:00:00Z", "indicator_types": ["malicious-activity"], "pattern_type": "sigma", "pattern": "title: Rare Scheduled Task Creations\nid: b20f6158-9438-41be-83da-a5a16ac90c2b\nstatus: experimental\ndescription: This rule detects rare scheduled task creations. Typically software gets installed on multiple systems and not only on a few. The aggregation and count\n    function selects tasks with rare names.\ntags:\n    - attack.persistence\n    - attack.t1053\n    - attack.s0111\nauthor: Florian Roth\ndate: 2017/03/17\nlogsource:\n    product: windows\n    service: taskscheduler\ndetection:\n    selection:\n        EventID: 106\n    timeframe: 7d\n    condition: selection | count() by TaskName < 5\nfalsepositives:\n    - Software installation\nlevel: low\n", "object_marking_refs": ["marking-definition--0a872d70-1d05-45dd-a25f-031af547a102"], "created_by_ref": "identity--082ca4f0-d825-45b6-aa42-ab9afd1e5503", "type": "indicator", "spec_version": "2.1", "id": "indicator--4401f7de-c808-42d3-8df1-383e4ce51848", "modified": "2020-04-02T10:14:31.795Z"}